# backend/deploy/systemd/mullusi-govern.service
# Purpose: manage the Mullusi Govern Cloud production Compose stack as a host service.
# Governance scope: startup ordering, controlled restart boundary, environment file witness, and deterministic stop command.
# Dependencies: systemd, Docker Engine with Compose v2, /opt/mullusi/govern-cloud deployment directory, and /etc/mullusi/govern.env.
# Invariants: service management never embeds secrets and never starts without the production environment file.

[Unit]
Description=Mullusi Govern Cloud API
Requires=docker.service
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/mullusi/govern-cloud
EnvironmentFile=/etc/mullusi/govern.env
ExecStart=/usr/bin/docker compose -f /opt/mullusi/govern-cloud/docker-compose.production.yaml up -d
ExecStop=/usr/bin/docker compose -f /opt/mullusi/govern-cloud/docker-compose.production.yaml down
TimeoutStartSec=120
TimeoutStopSec=60

[Install]
WantedBy=multi-user.target
